# templates/cron.d/chrysaor_nsupdate
# min	hr	date	month	wday	user	cmd
9,29,49	*	*	*	*	root	/bin/sh /root/cronjobs/nsupdate-script.sh
