# templates/cron.d/chrysaor_nsupdate
@hourly         root    /bin/sh /root/cronjobs/nsupdate-script.sh
